# 7.1/alpine3.8-nginx/rc.local
#!/bin/sh

/usr/sbin/nginx &

/usr/local/sbin/php-fpm &
